# deploy/storebot.service
[Unit]
Description=Storebot Telegram Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=storebot
WorkingDirectory=/opt/storebot
EnvironmentFile=/opt/storebot/.env
ExecStart=/opt/storebot/.venv/bin/storebot
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5
TimeoutStopSec=30
SyslogIdentifier=storebot

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/storebot/data /opt/storebot/backups
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
MemoryMax=512M

[Install]
WantedBy=multi-user.target
